# dot_anacron/cron.daily/executable_auto-rsync
#!/usr/bin/bash

pass show Remote/rsync | sshpass rsync -avzP --delete --delete-excluded --exclude-from="$HOME/.anacron/cron.daily/rsync-homedir-excludes.txt" /home mahid@192.168.1.6::NetBackup
